# roles/openslides-add-instance/templates/openslides_instance_proxy.socket.j2
[Socket]
ListenStream={{ openslides_instance_systemd_port_base + item|int }}

[Install]
WantedBy=sockets.target
