# roles/openslides-add-instance/templates/openslides_instance_redis.service.j2
[Unit]
Description=Openslides Instance
Wants=network.target

[Service]
ExecStart=/usr/bin/rkt run --insecure-options=image --port=6379-tcp:{{ openslides_instance_port_base }} docker://sameersbn/redis
NonBlocking=true
Restart=always
